# .docker/supervisor/supervisor.conf
[supervisord]
nodaemon=true

[program:app-chat]
command=php /application/bin/console chat:run
user=root
numprocs=1
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d